# hw12_13_14_15_calendar/.docker/pgsql/etc/profile.d/99-umask-regular-users.sh
#!/bin/sh

if [[ `id -u` -ge 1000 || `id -un` = "docker" ]]; then
    umask 0002
fi
